# crontab
*  *  *  *  *  bash /code/IDontWantToRead/startup.sh
0 */8  *  *  *  /usr/local/bin/php /code/IDontWantToRead/trigger_rescanFolders.php
