# crontab
*  *  *  *  *  bash /code/IDontWantToRead/startup.sh
*  *  *  *  *  echo "I can run" > /tmp/canIrun.txt
0 */8  *  *  *  /usr/local/bin/php /code/IDontWantToRead/trigger_rescanFolders.php
