# crontab
*  *  *  *  *    bash /code/IDontWantToRead/startup.sh
0 */8  *  *  *   php /code/IDontWantToRead/trigger_rescanFolders.php